# configs/supervisor_configs/gunicorn.conf
[program:flask_project]
command = gunicorn app:app -b 0.0.0.0:8000
directory = /home/ubuntu/twitter.rx/frontend
user = ubuntu
autostart=true
autorestart=true


